# 4 Extracted Root FS/_dcs-1130.dat.extracted/cramfs-root/etc/inittab
::sysinit:/etc/rc.d/rcS start
#::askfirst:/bin/ash
#tts/0::respawn:-/bin/sh
#tts/0::respawn:-/bin/ash
::respawn:-/bin/ash
#::ctrlaltdel:/etc/rc.d/rcK stop
#::shutdown:/etc/rc.d/rcK stop

